# systemd/switcher.service
[Unit]
Description=Proxy Switcher

[Service]
ExecStart=/usr/bin/switcher

[Install]
WantedBy=default.target
