# systemd/switcherd.service
[Unit]
Description=Proxy Switcher

[Service]
ExecStart=/usr/bin/switcherd

[Install]
WantedBy=default.target
